# _cocktails/250428_speedrun/23nation-2/src/app/app-shutdown.service
[Unit]
Description=...

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/home/ec2-user/shutdown.sh

[Install]
WantedBy=multi-user.target
