# _cocktails/250428_speedrun/23nation-2/src/app/crontab
* * * * * aws s3 cp /home/ec2-user/app/log/app.log "s3://wsi-logs-pmhn/accesslog-backup/interval/$(date +\%Y/\%m/\%d/\%H/\%M)/app.log"
